# net-print/pipslite/files/pipslited
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start () {
	ebegin "Starting pipslited"
	start-stop-daemon --quiet --start --exec /usr/bin/pipslited
	eend $?
}

stop () {
	ebegin "Stopping pipslited"
	start-stop-daemon --quiet --stop --exec /usr/bin/pipslited
	rm $(grep DummyDevicePath /etc/pipslite/pipslitedrc | sed -e "s/.*=//")
	eend $?
}
